# legacy/conf/supervisor/cron.conf
[program:cron]
command = runif RUN_CRON /usr/sbin/cron -f
priority = 100
autostart = true
