# conf/supervisor/nginx.conf
[program:nginx]
command = runif RUN_NGINX play /opt/ansible/nginx.yml nginx
priority = 600
autostart = true
startsecs = 10
environment = RUNIF_WAIT_SECS=10
