# conf/supervisor/sshd.conf
[program:sshd]
command = runif RUN_SSHD play /opt/ansible/openssh.yml root /usr/sbin/sshd -D
priority = 100
autostart = true
