# conf/supervisor/syslog-ng.conf
[program:syslog-ng]
command = runif RUN_SYSLOG /usr/sbin/syslog-ng --foreground
priority = 80
autostart = true
